# home/dot_config/systemd/user/atuin.service
[Unit]
Description=Atuin shell history daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=atuin daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
